# rootfs/services/clamd/run
#!/bin/bash
logger -p mail.info "s6-supervise : spawning clamd process"
exec clamd &>/dev/null
