# rootfs/services/clamd/run
#!/bin/bash
logger -p mail.info "s6-supervise : virus database downloaded, spawning clamd process"

echo $$>/var/run/clamav/clamd.pid
exec clamd &>/dev/null
