# rootfs/etc/cron.d/counters
00 00 * * * root echo 0 > /tmp/counters/clamd
00 00 * * * root echo 0 > /tmp/counters/dovecot
00 00 * * * root echo 0 > /tmp/counters/freshclam
00 00 * * * root echo 0 > /tmp/counters/postfix
00 00 * * * root echo 0 > /tmp/counters/rspamd
00 00 * * * root echo 0 > /tmp/counters/unbound
