# rootfs/etc/cron.d/spamassassin
30 02 * * * root /usr/bin/sa-learn --ham /var/mail/vhosts/*/*/mail/cur/* | logger -t sa-learn-ham
40 02 * * * root /usr/bin/sa-learn --spam /var/mail/vhosts/*/*/mail/.Spam/cur/* | logger -t sa-learn-spam
