# rootfs/services/rspamd/run
#!/bin/bash
logger -p mail.info "s6-supervise : spawning rspamd process"
exec rspamd -f -u _rspamd -g _rspamd
