# rootfs/services/dovecot/run
#!/bin/bash
logger -p mail.info "s6-supervise : spawning dovecot process"
exec dovecot -F &>/dev/null
