# rootfs/services/postfix/run
#!/bin/bash
logger -p mail.info "s6-supervise : spawning postfix process"
exec postfix -c /etc/postfix start-fg &>/dev/null
